# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:python]
command=python3 main.py

[program:serve]
command=serve ./frontend/dist 